# sdcard/debian/stuff/rc.local
#!/bin/sh -e

/usr/sbin/olinuxino_setup.sh
