# sdcard/debian/stuff/rc.local
#!/bin/sh -e

/usr/sbin/autoresize_partition.sh
